# pi-root/etc/systemd/system/ldrs-usb-hardening.service
[Unit]
Description=Sport Assist USB hardening (USBGuard + EEPROM)
DefaultDependencies=no
After=local-fs.target
Before=usbguard.service ldrs-wifi-ap.service
ConditionPathExists=/usr/local/bin/ldrs-apply-usb-hardening.sh

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/ldrs-apply-usb-hardening.sh

[Install]
WantedBy=multi-user.target
